# shairport-sync/rootfs/etc/services.d/shairport/run
#!/command/with-contenv sh
exec /usr/local/bin/run.sh
